# supervisor/symfony.conf
[program:symfony]
directory=/root/SYMFONYPROJECTNAME
command=/usr/bin/php app/console server:run 0.0.0.0:80
autostart=true
autorestart=true
redirect_stderr=true
